# luci-app-shadowtunnel/root/etc/uci-defaults/shadowtunnel
#!/bin/sh

chmod 777 /etc/shadowtunnel/*
chmod +x /etc/init.d/control_shadowt
exit 0
